# services/birdnet_sync.service
[Unit]
Description=Sync Service to backup and push Birdnet files to a server
[Service]
Type=simple
ExecStart=/home/phablet/Documents/.birdnet/services/bin/sync.sh
[Install]
WantedBy=default.target
