# services/birdnet_watchdog.timer
Description=Watchdog Service to monitor, log and repair
[Timer]
OnBootSec=6min
OnUnitActiveSec=11min
Unit=birdnet_watchdog.service
[Install]
WantedBy=default.target
