# services/birdnet_weekly_report.timer
[Unit]
Description=Weekly timer
[Timer]
OnCalendar=Weekly
Unit=birdnet_weekly_report.service
[Install]
WantedBy=default.target
